# rootfs/etc/s6-overlay/s6-rc.d/init-services/up
# This S6 service does not do anything. It is used to signal the start of initializing container
# longrun services to maintain a proper order of dependencies in child images.