# rootfs/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/command/execlineb -P
exec /usr/sbin/cron -f -L 5
